# debian/etc/udev/rules.d/10-usb-storage.rules
# do not edit this file, it will be overwritten on update
KERNEL!="sd[a-z][0-9]",GOTO="automount_exit"
ACTION=="add", SUBSYSTEM=="block", RUN+="/bin/mkdir -p /media/%k", RUN+="/usr/bin/mount $devnode /media/%k"
ACTION=="remove", RUN+="/usr/bin/umount /media/%k",RUN+="/bin/rmdir /media/%k"
LABEL="automount_exit"

