# debian/etc/udev/rules.d/12-sdcard.rules
# do not edit this file, it will be overwritten on update
KERNEL!="mmcblk[1-2]p[1-9]",GOTO="automount_exit"
ACTION=="add", SUBSYSTEM=="block",RUN+="/bin/mkdir -p /media/sdcard%n",RUN+="/usr/bin/mount $devnode /media/sdcard%n"
ACTION=="remove", RUN+="/usr/bin/systemd-umount /media/sdcard%n",RUN+="/bin/rmdir /media/sdcard%n"
LABEL="automount_exit"
