# debian/etc/udev/rules.d/70-isp.rules
# do not edit this file, it will be overwritten on update
ACTION=="add", SUBSYSTEM=="platform", DRIVER=="hobot,x3-isp", MODE="0666", RUN+="/bin/chmod 666 /sys/devices/platform/soc/b3000000.isp/v3a"

