# debian/etc/udev/rules.d/99-hobot-group-management.rules
# do not edit this file, it will be overwritten on update

ACTION=="add", KERNEL=="bpu*", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="ac_*", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="rtc*", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="ubi_ctrl", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="log_*", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="console*", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="ion", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="port_*", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="hmipc", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"
ACTION=="add", KERNEL=="sensor_ctrl", GROUP="misc", OWNER="misc", MODE="0664", RUN+="/bin/chown :misc %k"

ACTION=="add", KERNEL=="i2c*", GROUP="i2c", OWNER="i2c", MODE="0664", RUN+="/bin/chown :i2c %k"
ACTION=="add", KERNEL=="gpiochip*", GROUP="gpio", OWNER="gpio", MODE="0664", RUN+="/bin/chown :gpio %k"

ACTION=="add", KERNEL=="vpu*", GROUP="vpu", OWNER="vpu", MODE="0664", RUN+="/bin/chown :vpu %k"
ACTION=="add", KERNEL=="jpu*", GROUP="jpu", OWNER="jpu", MODE="0664", RUN+="/bin/chown :jpu %k"

ACTION=="add", KERNEL=="sif*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="mipi*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="ipu*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="isp*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="pym*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="ldc*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="dwe*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="vps_*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="gdc*", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="osd", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="vio_bind_info", KERNEL=="osd", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"
ACTION=="add", KERNEL=="vio_mp", KERNEL=="osd", GROUP="vps", OWNER="vps", MODE="0664", RUN+="/bin/chown :vps %k"

ACTION=="add", KERNEL=="iar*", GROUP="graphics", OWNER="graphics", MODE="0664", RUN+="/bin/chown :graphics %k"

ACTION=="add", KERNEL=="video*", GROUP="video", OWNER="video", MODE="0664", RUN+="/bin/chown :video %k"
ACTION=="add", KERNEL=="v4l*", GROUP="video", OWNER="video", MODE="0664", RUN+="/bin/chown :video %k"

SUBSYSTEM=="sound", GROUP="misc", OWNER="sound", OPTIONS+="static_node=snd/seq", OPTIONS+="static_node=snd/timer"
ACTION=="add", SUBSYSTEM=="sound", GROUP="misc", OWNER="sound", MODE="0664", RUN+="/bin/chown :misc %k"

SUBSYSTEM=="gpio", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R $env{SUBSYSTEM} %N %S/$env{DEVPATH} %S/class/misc;chmod g+w %N %S/$env{DEVPATH} %S/class/gpio'"
SUBSYSTEM=="vps", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R $env{SUBSYSTEM} %N %S/$env{DEVPATH} %S/class/vps;chmod g+rw -R %N %S/$env{DEVPATH} %S/class/vps'"
SUBSYSTEM=="ipu", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R $env{SUBSYSTEM} %N %S/$env{DEVPATH} ;chmod g+w %N %S/$env{DEVPATH}'"
SUBSYSTEM=="misc", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R $env{SUBSYSTEM} %N %S/$env{DEVPATH} %S/class/misc;chmod g+w -R %N %S/$env{DEVPATH} %S/class/misc'"
SUBSYSTEM=="jpu", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R $env{SUBSYSTEM} %N %S/$env{DEVPATH} %S/class/jpu;chmod g+w -R %N %S/$env{DEVPATH} %S/class/jpu'"
SUBSYSTEM=="graphics", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R $env{SUBSYSTEM} %N %S/$env{DEVPATH} %S/class/graphics;chmod g+w -R %N %S/$env{DEVPATH} %S/class/graphics'"
SUBSYSTEM=="vpu", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R $env{SUBSYSTEM} %N %S/$env{DEVPATH} %S/class/vpu;chmod g+w -R %N %S/$env{DEVPATH} %S/class/vpu'"
SUBSYSTEM=="vpu_clk_dev", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp -R vpu %N %S/$env{DEVPATH} %S/class/vpu_clk_dev;chmod g+w -R %N %S/$env{DEVPATH} %S/class/vpu_clk_dev'"
SUBSYSTEM=="platform", DRIVER=="hb_vpu", ACTION=="add", \
        RUN+="/bin/bash -c 'chgrp vpu %S/devices/platform/soc/a8000000.vpu/fps;chmod g+w %S/devices/platform/soc/a8000000.vpu/fps'"

SUBSYSTEM=="mem", KERNEL=="mem", MODE="0660"
