# debian/etc/systemd/system/systemd-udevd.service.d/delay_load_usb_uvc_driver.conf
[Unit]
After=systemd-sysusers.service systemd-hwdb-update.service systemd-modules-load.service
Requires=systemd-modules-load.service

