# debian/DEBIAN/prerm
#!/bin/bash

update-rc.d usbsel remove
update-rc.d hobot-rc remove

systemctl disable usbsel.service
systemctl disable hobot-rc.service

exit 0

